# zurmo_elasticsearch/fleet/zurmo_elasticsearch_discovery@.service
[Unit]
Description=Announces Elasticsearch
BindsTo=zurmo_elasticsearch@%i.service
After=zurmo_elasticsearch@%i.service

[Service]
EnvironmentFile=/etc/environment

ExecStart=/bin/sh -c "UUID=`uuidgen`; echo $UUID > /tmp/zurmo_elasticsearch@%i.uuid; TTL=10; SLEEP_TIME=7 \
  etcdctl setdir /services/logstorage/$UUID --ttl $TTL; \
  etcdctl setdir /services/logstorage/$UUID/settings --ttl $TTL; \
    while true; \
      do \
        etcdctl updatedir /services/logstorage/$UUID --ttl $TTL; \
        etcdctl set /services/logstorage/$UUID/host %H --ttl $TTL; \
        etcdctl set /services/logstorage/$UUID/ip ${COREOS_PRIVATE_IPV4} --ttl $TTL; \
        etcdctl set /services/logstorage/$UUID/http-port 9200 --ttl $TTL; \
        etcdctl set /services/logstorage/$UUID/transport-port 9300 --ttl $TTL; \
        sleep $SLEEP_TIME; \
      done"

ExecStop=/bin/sh -c "etcdctl rm --recursive /services/logstorage/`cat /tmp/zurmo_elasticsearch@%i.uuid` && rm -f /tmp/zurmo_elasticsearch%i.uuid"

[X-Fleet]
MachineOf=zurmo_elasticsearch@%i.service